# recipes-geraldpi/sync-scripts/files/sync-scripts.timer
[Unit]
Description=Periodic Sync docs between all geraldpi-document instances

[Timer]
OnCalendar=*-*-* 02:00:00
RandomizedDelaySec=600

[Install]
WantedBy=timers.target

