# recipes-geraldpi/thermostat/files/temp-scanner.service
[Unit]
Description=Temp Scanner Service

[Service]
Type=simple
ExecStart=/opt/thermostat/scan.py
Restart=Always

[Install]
WantedBy=multi-user.target

