# recipes-geraldpi/thermostat/files/thermostat.service
[Unit]
Description=Thermostat Run Script
#After=network.target

[Service]
Type=oneshot
ExecStart=/opt/thermostat/run.sh

